# scripts/systemd/vps-agent.service
[Unit]
Description=VPS Monitor Agent
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/vps-monitor
ExecStart=/opt/vps-monitor/bin/vps-agent -config /opt/vps-monitor/config/agent.yaml
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target